# deploy/systemd/jarvis-live.service
[Unit]
Description=JARVIS live context + supervisor loop (Evolutionary Trading Algo)
Documentation=file://__INSTALL_DIR__/deploy/README.md
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify lets the daemon report READY=1 + WATCHDOG=1 via sd_notify.
# The daemon implementation lives in eta_engine.obs.watchdog --
# see WatchdogPinger() and notify_ready().
Type=notify
NotifyAccess=main
WatchdogSec=60
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=__INSTALL_DIR__/.venv/bin/python -m eta_engine.scripts.jarvis_live \
  --inputs __INSTALL_DIR__/docs/premarket_inputs.json \
  --out-dir %h/.local/state/eta_engine \
  --interval 60
Restart=always
RestartSec=15
# If WatchdogSec triggers, escalate to a full restart (don't just kill).
WatchdogSignal=SIGTERM
StandardOutput=append:%h/.local/log/eta_engine/jarvis-live.log
StandardError=append:%h/.local/log/eta_engine/jarvis-live.err

# ------------------------------------------------------------------
# Hardening (see obs.vps_hardening for rationale)
# ------------------------------------------------------------------
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/eta_engine %h/.local/log/eta_engine __INSTALL_DIR__/docs __INSTALL_DIR__/runs
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

[Install]
WantedBy=default.target
